# board/freescale/imx6ull/rootfs-overlay/etc/inittab
::sysinit:/etc/init.d/rcS
ttymxc0::respawn:/bin/sh
::shutdown:/bin/umount -a -r
